# systemd-units/timer-afterboot.timer
[Unit]
Description=Timer for event after successfull boot

[Timer]
OnActiveSec=90s
Unit=timer-afterboot.target

[Install]
WantedBy=multi-user.target
